# client/frontend/misc/monitor-ui@.service
[Unit]
Description=monitor sysinfo

[Service]
Slice=session.slice
Type=notify
ExecStart=/usr/bin/monitor-server -a %i